# src/systemd/slowcontrol.service.in
# file installed from ${CMAKE_CURRENT_SOURCE_DIR}
[Unit]
Description=Slowcontrol ${DAEMON} start
After=network.target nss-lookup.target time-sync.target

[Service]
Type=forking
PrivateTmp=false
PIDFile=/var/run/${DAEMON}.pid
Environment="HOME=/root"
ExecStart=/bin/bash -c '${CMAKE_INSTALL_PREFIX}/bin/${DAEMON} --pidFile /var/run/${DAEMON}.pid &> /tmp/${DAEMON}.log'
KillMode=mixed
Restart=always
${SERVICE_EXTRALINES}

[Install]
WantedBy=${WANTED_BY}
