# src/systemd/slowcontrol.service.in
# file installed from ${CMAKE_CURRENT_SOURCE_DIR}
[Unit]
Description=Slowcontrol ${DAEMON} start
After=network.target nss-lookup.target time-sync.target
${UNIT_EXTRALINES}

[Service]
Type=simple
PrivateTmp=false
PIDFile=/run/${DAEMON}.pid
Environment="HOME=/root"
ExecStart=${CMAKE_INSTALL_PREFIX}/bin/${DAEMON} --nodaemon --pidFile /var/run/${DAEMON}.pid
Restart=always
${SERVICE_EXTRALINES}

[Install]
WantedBy=${WANTED_BY}
